# backuper.timer
[Unit]
Description=Run backuper every Saturday

[Timer]
OnCalendar=Sat *-*-* 12:00:00
Persistent=true

[Install]
WantedBy=timers.target
